# images/twemproxy/files/supervisord.conf
[supervisord]
nodaemon=true

[program:nutcracker]
command=nutcracker --conf-file=/etc/twemproxy/twemproxy.yaml
autorestart=true
